# unit-files/helloworld/app@.service.erb
[Unit]
After=install-%p.service
After=akka-cluster-nodes.service
After=docker.service
After=weave.service
Description=CoreOS Akka Cluster Application
Requires=install-%p.service
Requires=akka-cluster-nodes.service
Requires=docker.service
Requires=weave.service

[Service]
User=core
TimeoutStartSec=0
KillMode=none

# Read the local environment (needed to to get COREOS_PRIVATE_IPV4)
EnvironmentFile=/etc/environment

# Remove old docker containers with the same name if they exist
ExecStartPre=-/usr/bin/docker kill %p-%i
ExecStartPre=-/usr/bin/docker rm %p-%i

# Start the docker container with weave networking (we pass in the etcd URL via the environment)
# NOTE: we need to use an IP address here as Akka will otherwise map the DNS name to the wrong NIC address!
ExecStartPre=/usr/bin/sudo /opt/bin/weave run 10.42.0.%i/16 --name %p-%i -h %p-%i.weave.local \
  -e "APP_ADDR=10.42.0.%i" \
  -e "ETCD_URL=http://${COREOS_PRIVATE_IPV4}:4001" \
  <%= @docker[:app] %>
ExecStart=/usr/bin/docker attach %p-%i

# When stopped, stop the docker container
ExecStop=/usr/bin/docker stop %p-%i

[X-Fleet]
# Only run cluster nodes on Akka equipped compute nodes
MachineMetadata=type=akka
