# unit-files/seed@.service.erb
[Unit]
After=akka-cluster-nodes.service
Description=Used to form a cluster from an initial seed node
Requires=akka-cluster-nodes.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutSec=0

# Ensure seed node exists
ExecStartPre=/usr/bin/etcdctl get /akka.cluster.nodes/%i
ExecStart=/usr/bin/etcdctl set /akka.cluster.nodes/%i MemberUp
