# unit-files/seed@.service.erb
[Unit]
After=etcd.service
Description=Used to form a cluster from an initial seed node
Requires=etcd.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutSec=0

# Ensure seed node exists
ExecStartPre=/usr/bin/etcdctl get /<%= @service_discovery[:akka] %>/%i
ExecStart=/usr/bin/etcdctl set /<%= @service_discovery[:akka] %>/%i Up
