# unit-files/vulcand.service.erb
<%
  @image = @docker[:vulcand]
  @logLevel = "INFO" # INFO, WARN or ERROR
  @ip_address = "10.42.253.253"
  @port = 8080
  @apiPort = 8182
%>

[Unit]
Description=Vulcand Load Balancer
After=docker.service
After=weave.service
Requires=docker.service
Requires=weave.service

[Service]
User=core
TimeoutStartSec=0

# Read the local environment (needed to to get COREOS_PRIVATE_IPV4)
EnvironmentFile=/etc/environment

ExecStartPre=-/usr/bin/docker kill %p
ExecStartPre=-/usr/bin/docker rm %p
ExecStartPre=/usr/bin/docker pull <%= @image %>

ExecStartPre=/usr/bin/sudo /opt/bin/weave run <%= @ip_address %>/16 --name %p -h %p.weave.local \
  -p <%= @port %>:<%= @port %> \
  -p <%= @apiPort %>:<%= @apiPort %> \
  <%= @image %> \
  /go/bin/vulcand -apiInterface=0.0.0.0 -interface=0.0.0.0 -logSeverity=<%= @logLevel %> -etcd=http://${COREOS_PRIVATE_IPV4}:4001 -port=<%= @port %> -apiPort=<%= @apiPort %>
ExecStart=/usr/bin/docker attach %p

ExecStop=/usr/bin/docker stop %p

[X-Fleet]
MachineMetadata=load-balancer=true
